# scripts/com.philly-fair-measure.snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.philly-fair-measure.snapshot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nhand/.local/bin/uv</string>
        <string>run</string>
        <string>philly</string>
        <string>snapshot-all</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/nhand/DataProjects/philly-fair-measure</string>
    <!-- Weekly: Mondays 07:00 local. launchd runs missed jobs at next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/nhand/DataProjects/philly-fair-measure/data/logs/snapshot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/nhand/DataProjects/philly-fair-measure/data/logs/snapshot.log</string>
</dict>
</plist>
